# 99-xbelite2.rules
# xbelite2 — Xbox Elite Series 2 userspace driver
#
# Automatically bind Elite 2 to hid-xbelite2 instead of hid-microsoft.

# BT BLE (0x0B22): unbind whatever grabbed it, bind xbelite2
ACTION=="add", SUBSYSTEM=="hid", KERNEL=="0005:045E:0B22.*", \
    RUN+="/bin/sh -c 'echo %k > /sys/bus/hid/drivers/microsoft/unbind 2>/dev/null; echo %k > /sys/bus/hid/drivers/xbelite2/bind 2>/dev/null; true'"

# BT Classic (0x0B05)
ACTION=="add", SUBSYSTEM=="hid", KERNEL=="0005:045E:0B05.*", \
    RUN+="/bin/sh -c 'echo %k > /sys/bus/hid/drivers/microsoft/unbind 2>/dev/null; echo %k > /sys/bus/hid/drivers/xbelite2/bind 2>/dev/null; true'"

# USB (0x0B00): unbind xpad
ACTION=="add", SUBSYSTEM=="hid", KERNEL=="0003:045E:0B00.*", \
    RUN+="/bin/sh -c 'echo %k > /sys/bus/hid/drivers/xpad/unbind 2>/dev/null; echo %k > /sys/bus/hid/drivers/xbelite2/bind 2>/dev/null; true'"

# hidraw: ROOT ONLY — prevents Steam from seeing the raw device.
# Only our daemon (running as root) accesses hidraw.
# Steam/games see the virtual gamepad at /dev/input/event* instead.
KERNEL=="hidraw*", ATTRS{idVendor}=="045e", ATTRS{idProduct}=="0b22", MODE="0600", GROUP="root", TAG-="uaccess"
KERNEL=="hidraw*", ATTRS{idVendor}=="045e", ATTRS{idProduct}=="0b05", MODE="0600", GROUP="root", TAG-="uaccess"
KERNEL=="hidraw*", ATTRS{idVendor}=="045e", ATTRS{idProduct}=="0b00", MODE="0600", GROUP="root", TAG-="uaccess"

# uinput access for the daemon
KERNEL=="uinput", MODE="0660", GROUP="input"
