# pkg/xbelite2/etc/udev/rules.d/99-xbelite2.rules
# xbelite2 - Xbox Elite Series 2 userspace driver
#
# Prevent xpad/xpadneo from claiming Elite 2 controllers so our daemon
# can access them via hidraw.

# BT: grant hidraw access
KERNEL=="hidraw*", ATTRS{idVendor}=="045e", ATTRS{idProduct}=="0b05", MODE="0660", TAG+="uaccess"
KERNEL=="hidraw*", ATTRS{idVendor}=="045e", ATTRS{idProduct}=="0b22", MODE="0660", TAG+="uaccess"
KERNEL=="hidraw*", ATTRS{idVendor}=="045e", ATTRS{idProduct}=="028e", SUBSYSTEM=="hidraw", MODE="0660", TAG+="uaccess"

# USB: unbind xpad/xpadneo and rebind to usbhid for hidraw access
ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="045e", ATTR{idProduct}=="0b00", \
    RUN+="/bin/sh -c 'for drv in xpad xpadneo; do echo $kernel > /sys/bus/usb/drivers/$drv/unbind 2>/dev/null; done; echo $kernel > /sys/bus/usb/drivers/usbhid/bind 2>/dev/null'"

# BT: unbind xpadneo from Elite 2 BLE devices so hid-generic handles them
ACTION=="add", SUBSYSTEM=="hid", ATTRS{idVendor}=="045e", ATTRS{idProduct}=="0b22", \
    RUN+="/bin/sh -c 'echo $kernel > /sys/bus/hid/drivers/xpadneo/unbind 2>/dev/null'"
ACTION=="add", SUBSYSTEM=="hid", ATTRS{idVendor}=="045e", ATTRS{idProduct}=="0b05", \
    RUN+="/bin/sh -c 'echo $kernel > /sys/bus/hid/drivers/xpadneo/unbind 2>/dev/null'"

# USB: grant hidraw access after rebind
KERNEL=="hidraw*", ATTRS{idVendor}=="045e", ATTRS{idProduct}=="0b00", MODE="0660", TAG+="uaccess"

# uinput access for virtual gamepad creation
KERNEL=="uinput", MODE="0660", GROUP="input"
