# pkg/xbelite2/etc/udev/rules.d/99-xbelite2.rules
# xbelite2 — Xbox Elite Series 2 userspace driver
#
# Automatically unbind hid_microsoft from Elite 2 controllers after probe.
# This prevents duplicate input devices while allowing hid_microsoft to
# remain loaded for other Microsoft devices.

# BT BLE (0x0B22): unbind hid_microsoft, keep xbelite2
ACTION=="add", SUBSYSTEM=="hid", ATTRS{idVendor}=="045e", ATTRS{idProduct}=="0b22", \
    RUN+="/bin/sh -c 'sleep 0.1; for d in /sys/bus/hid/drivers/microsoft/0005:045E:0B22.*; do [ -e \"$$d\" ] && echo $$(basename $$d) > /sys/bus/hid/drivers/microsoft/unbind 2>/dev/null; done; true'"

# BT Classic (0x0B05): unbind hid_microsoft, keep xbelite2  
ACTION=="add", SUBSYSTEM=="hid", ATTRS{idVendor}=="045e", ATTRS{idProduct}=="0b05", \
    RUN+="/bin/sh -c 'sleep 0.1; for d in /sys/bus/hid/drivers/microsoft/0005:045E:0B05.*; do [ -e \"$$d\" ] && echo $$(basename $$d) > /sys/bus/hid/drivers/microsoft/unbind 2>/dev/null; done; true'"

# USB (0x0B00): prevent xpad if it binds
ACTION=="bind", SUBSYSTEM=="usb", DRIVER=="xpad", ATTRS{idVendor}=="045e", ATTRS{idProduct}=="0b00", \
    RUN+="/bin/sh -c 'echo %k > /sys/bus/usb/drivers/xpad/unbind 2>/dev/null; true'"

# Misc char devices created by xbelite2 kmod — grant input group r/w so the
# GUI / xbe2-rw can issue profile writes without sudo.
KERNEL=="xbelite2",    SUBSYSTEM=="misc", MODE="0660", GROUP="input"
KERNEL=="xbelite2_bt", SUBSYSTEM=="misc", MODE="0660", GROUP="input"
