# pkg/xbelite2.install
post_install() {
  udevadm control --reload-rules
  udevadm trigger
  systemctl daemon-reload
  echo "------------------------------------------------------"
  echo "To start and enable the service, run:"
  echo "  sudo systemctl enable --now xbelite2d.service"
  echo "------------------------------------------------------"
}

post_upgrade() {
  echo "-----------------------------------------------------------------"
  echo "If you have the service running, don't forget to restart it:"
  echo "  sudo systemctl restart xbelite2d.service"
  echo "-----------------------------------------------------------------"
}

pre_remove() {
  systemctl stop xbelite2d.service 2>/dev/null
}
